# docker-deps/supervisord.conf
[supervisord]
nodaemon=true

[inet_http_server]
port: *:8080

[program:php-webserver]
command=php -S 0.0.0.0:80 -t %(ENV_APP_PATH)s/web/
user=root
